# systemctl/hue.service
[Unit]
Description=NZXT Hue Plus Init

[Service]
ExecStart=%h/.local/bin/hue profile apply bluepink

[Install]
WantedBy=default.target
